# conf/supervisor/conf.d/app.conf
[program:uwsgi-app]
command = /usr/local/bin/uwsgi --ini /etc/uwsgi/magiclantern.ini

[program:nginx-app]
command = /usr/sbin/nginx -c /etc/nginx/nginx.conf
