# conf/supervisor/vis.conf
[program:uwsgi-vis]
command = /usr/local/bin/uwsgi --ini /etc/uwsgi/vis.ini

[program:nginx-vis]
command = /usr/sbin/nginx -c /etc/nginx/nginx.conf

[program:confd-vis]
command = /usr/local/bin/confd --backend=env
